# scripts/pm-listen.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://purl.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.elmspark.pm-listen</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/kennjordan/Developer/elmspark/mac-mini-agent/scripts/pm-listen-service.sh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/tmp/steer/listen-service.out</string>
	<key>StandardErrorPath</key>
	<string>/tmp/steer/listen-service.err</string>
	<key>WorkingDirectory</key>
	<string>/Users/kennjordan/Developer/elmspark/mac-mini-agent</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
	</dict>
</dict>
</plist>
